# deploy/templates/dna-backup.cron.j2
0 */6 * * * www-data cp /opt/dna/pbp.db /opt/dna/backups/pbp-$(date +\%Y\%m\%d-\%H\%M).db
